# supervisord.conf
; Self-contained supervisord config. Overwrites the stock top-level config;
; we do not rely on conf.d/*.conf include behavior.

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ---- programs ----

[program:xvfb]
priority=100
command=/usr/bin/Xvfb :0 -screen 0 1440x900x24 -dpi 96 -nolisten tcp
autorestart=true
startretries=3
startsecs=2
environment=HOME="/root"
stdout_logfile=/tmp/xvfb.out
stderr_logfile=/tmp/xvfb.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:dbus]
priority=200
command=/usr/bin/dbus-daemon --session --nofork --address=unix:path=/tmp/dbus-session-bus
autorestart=true
startretries=3
startsecs=2
environment=HOME="/root"
stdout_logfile=/tmp/dbus.out
stderr_logfile=/tmp/dbus.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:fluxbox]
priority=300
command=/usr/bin/fluxbox
autorestart=true
startretries=3
startsecs=2
environment=HOME="/root",DISPLAY=":0"
stdout_logfile=/tmp/fluxbox.out
stderr_logfile=/tmp/fluxbox.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:x11vnc]
priority=400
command=/usr/bin/x11vnc -display :0 -forever -shared -rfbport 5901 -nopw -xkb -noxdamage
autorestart=true
startretries=3
startsecs=2
environment=HOME="/root"
stdout_logfile=/tmp/x11vnc.out
stderr_logfile=/tmp/x11vnc.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:claude-desktop]
priority=500
command=/usr/local/bin/start-claude.sh
autorestart=true
startretries=3
startsecs=2
environment=HOME="/root",DISPLAY=":0",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session-bus"
stdout_logfile=/tmp/claude-desktop.out
stderr_logfile=/tmp/claude-desktop.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
